# systemd/user/envisaged-web-dev.service
[Unit]
Description=Envisaged Web UI (development, hot reload)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/utensils/Envisaged
ExecStart=/run/current-system/sw/bin/bash -lc 'cd %h/Projects/utensils/Envisaged && ./scripts/web-dev'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
